# deploy/systemd/slack-bridge.service
[Unit]
Description=PeerMarket Slack bridge (socket mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=peermarket-agent
Group=peermarket-agent
WorkingDirectory=/opt/peermarket-agent
EnvironmentFile=/etc/peermarket-agent/secrets.env
ExecStartPre=/bin/sh -c 'install -d -o peermarket-agent -g peermarket-agent -m 0700 "${VIDEO_MEDIA_ROOT:-/var/peermarket-agent/video-media}"'
ExecStartPre=/usr/bin/ffmpeg -version
ExecStartPre=/usr/bin/ffprobe -version
ExecStart=/opt/peermarket-agent/.venv/bin/peermarket-slack-bridge
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/peermarket-agent
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
